# misc/ras-mc-ctl.service.in
[Unit]
Description=Initialize EDAC v3.0.0 Drivers For Machine Hardware

[Service]
Type=oneshot
ExecStart=@sbindir@/ras-mc-ctl --register-labels
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
